# deploy/systemd/checkbox-shift-scheduler.service
# NOTE: this is a production template
# DO NOT APPLY DIRECTLY WITHOUT REVIEW

[Unit]
Description=Checkbox Shift Scheduler Service
After=network.target

[Service]
User=root
WorkingDirectory=/root/inventory/app/services
Environment="PYTHONPATH=/root/inventory"
EnvironmentFile=/root/inventory/.env
ExecStart=/root/inventory/.venv/bin/python -m app.services.checkbox_shift_scheduler_service
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
