# deploy/systemd/competitor_price_scheduler.service
# NOTE: this is production snapshot from server
# DO NOT APPLY DIRECTLY WITHOUT REVIEW

[Unit]
Description=Competitor Price Scheduler Service
After=network.target

[Service]
User=root
WorkingDirectory=/root/inventory/app/services
Environment="PYTHONPATH=/root/inventory"
EnvironmentFile=/root/inventory/.env
# normalized to python -m for consistent runtime behavior
ExecStart=/root/inventory/.venv/bin/python -m app.services.competitor_price_scheduler
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
