# deploy/systemd/telegram_bot.service
# NOTE: this is production snapshot from server
# DO NOT APPLY DIRECTLY WITHOUT REVIEW
# WARNING: active but disabled on server

[Unit]
Description=Telegram Bot Service
After=network.target

[Service]
User=root
WorkingDirectory=/root/inventory/app/services
Environment="PYTHONPATH=/root/inventory"
EnvironmentFile=/root/inventory/.env
# normalized to python -m for consistent runtime behavior
ExecStart=/root/inventory/.venv/bin/python -m app.services.telegram_bot
Restart=always

[Install]
WantedBy=multi-user.target
